# supervisord.conf
[supervisord]
logfile=supervisord.log
pidfile=supervisord.pid

[program:telegram_bot]
command=python bot.py
autostart=true
autorestart=true
stderr_logfile=bot.err.log
stdout_logfile=bot.out.log
